# templates/humio@.service.j2
[Unit]
  Description=Humio service
  After=network.service

[Service]
  Type=notify
  User=humio
  Group=humio
  LimitNOFILE=250000:250000
  EnvironmentFile=/etc/humio/server_%i.conf
  EnvironmentFile=/etc/humio/server_all.conf
  EnvironmentFile=/etc/humio/server_user_%i.conf
  WorkingDirectory={{ humio_data_path }}/{{ humio_host_id }}-%i
  ExecStart=/usr/bin/numactl --cpunodebind=%i --membind=%i /usr/bin/java -XX:+UnlockExperimentalVMOptions -XX:+UseZGC -Xlog:gc -Xss2M -XX:MaxDirectMemorySize=32G -Dhumio.auditlog.dir=${HUMIO_AUDITLOG_DIR} -Dhumio.debuglog.dir=${HUMIO_DEBUGLOG_DIR} -jar /usr/lib/humio/server-{{ humio_version }}.jar

[Install]
  WantedBy=default.target
